# telegram-bot/telegram-bot.service
[Unit]
Description=Telegram Bot for Kamila Zolotova
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/cybersec-psychology/telegram-bot
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/cybersec-psychology/telegram-bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target